# ansible/roles/container-backup/templates/docker-backup.timer.j2
[Unit]
Description=Nightly backup for Docker project {{ docker_project_name }}

[Timer]
OnCalendar=*-*-* 03:00:00
Persistent=true

[Install]
WantedBy=timers.target
